# system/pkg/adaptrandr/99-adaptrandr.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="echo"
